# images/supervisord-images.conf
[program:hivewing-images]
command=/home/git/start-images.sh
user=git
numprocs=1
autostart=true
autorestart=true
