# images/supervisord.conf
[supervisord]

[program:hivewing]
command=/home/git/start-images.sh
directory=/home/git
user=git
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/home/git/hivewing-images.log
