# web/supervisord.conf
[supervisord]

[program:hivewing]
command=/home/hivewing/start-web.sh
directory=/home/hivewing
user=hivewing
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/home/hivewing/logs/hivewing-web.log
